# scripts/com.qplus.surgery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.2.dtd">
<!--
  macOS launchd 서비스 자동 실행 설정 파일
  경로: ~/Library/LaunchAgents/com.qplus.surgery.plist
  
  사용 방법:
  1. 본 파일을 ~/Library/LaunchAgents/com.qplus.surgery.plist 에 복사합니다.
  2. <string> 안의 경로들을 실제 맥미니 환경에 맞춰 수정합니다. (예: /Users/사용자명/Developer/lgbtqplus_surgery)
  3. 실제 앱 설정은 레포 루트의 .env 파일에서 읽습니다.
  4. 다음 명령어로 서비스를 등록하고 시작합니다:
     launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.qplus.surgery.plist
     launchctl enable gui/$(id -u)/com.qplus.surgery
     launchctl kickstart -k gui/$(id -u)/com.qplus.surgery
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qplus.surgery</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jsbang/Developer/lgbtqplus_surgery/.venv/bin/uvicorn</string>
        <string>api_server:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/jsbang/Developer/lgbtqplus_surgery</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/jsbang/Developer/lgbtqplus_surgery/scripts/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jsbang/Developer/lgbtqplus_surgery/scripts/stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin</string>
    </dict>
</dict>
</plist>
